# docker/php/supervisord.prod.conf
[supervisord]
nodaemon=true
user=laravel
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord/supervisord.pid
loglevel=info

[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
startretries=5
user=laravel
redirect_stderr=true
stdout_logfile=/var/log/supervisor/php-fpm.out.log
stderr_logfile=/var/log/supervisor/php-fpm.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=10

[program:horizon]
command=php /var/www/html/artisan horizon
directory=/var/www/html
autostart=true
autorestart=true
startretries=5
user=laravel
redirect_stderr=true
stdout_logfile=/var/log/supervisor/horizon.out.log
stderr_logfile=/var/log/supervisor/horizon.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=20
stopsignal=TERM
stopwaitsecs=30

[program:schedule]
command=php /var/www/html/artisan schedule:work
directory=/var/www/html
autostart=true
autorestart=true
startretries=5
user=laravel
redirect_stderr=true
stdout_logfile=/var/log/supervisor/schedule.out.log
stderr_logfile=/var/log/supervisor/schedule.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=30
